# debian/pinephone-tweaks.prerm
#!/bin/sh

if [ "$1" = "remove" ]; then
    # Disable automatic camera setup
    systemctl disable pinephone-camera-setup.service
fi
